# bin/oms-web.service
[Unit]
Description=oms-web

[Service]
WorkingDirectory=/opt/oms
PrivateTmp=true
Restart=always
Type=simple
ExecStart=/usr/java/jdk1.8.0_181-amd64/jre/bin/java -jar -Xmx4g -Xms1g -Xss512k -Dloader.path=. -Dspring.profiles.active=web-prod /opt/oms/lib/oms-web-3.0-SNAPSHOT.jar

[Install]
WantedBy=multi-user.target
